# systemd/mcp-context-browser.service
# MCP Context Browser - User Service Unit
#
# This is a systemd user service file. Install to:
#   ~/.config/systemd/user/mcp-context-browser.service
#
# Commands:
#   systemctl --user daemon-reload
#   systemctl --user enable mcp-context-browser
#   systemctl --user start mcp-context-browser
#   systemctl --user status mcp-context-browser
#   journalctl --user -u mcp-context-browser -f

[Unit]
Description=MCP Context Browser - Semantic Code Search Server
Documentation=https://github.com/marlonsc/mcp-context-browser
After=default.target

[Service]
Type=simple

# Binary location (use %h for $HOME expansion)
ExecStart=%h/.local/bin/mcp-context-browser --config %h/.config/mcp-context-browser/config.toml

# Reload via SIGUSR1 triggers exec-based respawn (preserves PID)
ExecReload=/bin/kill -SIGUSR1 $MAINPID

# Graceful shutdown configuration
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Auto-restart on failure
Restart=on-failure
RestartSec=5

# Environment configuration
Environment=RUST_LOG=info
Environment=MCP_TRANSPORT_MODE=hybrid
Environment=MCP_PORT=3001

# Working directory for data storage
WorkingDirectory=%h/.local/share/mcp-context-browser

# Resource limits (optional, adjust as needed)
# MemoryMax=2G
# CPUQuota=100%

[Install]
WantedBy=default.target
